# airootfs/etc/profile.d/stack.sh
#!/bin/bash

PATH="/usr/local/stack:${PATH}"
